# data_gatherer/deploy/polymarket-archiver.service
[Unit]
Description=Polymarket Weather Archiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-weather-trading/data_gatherer
Environment=CONFIG=/opt/polymarket-weather-trading/data_gatherer/config.json
Environment=INTERVAL_S=60
Environment=DISCOVER_EVERY_S=3600
Environment=LOG_FILE=/var/log/polymarket-archiver/archiver.log
ExecStart=/usr/bin/env bash /opt/polymarket-weather-trading/data_gatherer/run_archiver.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

